# kiosk-configuration/ensure-kiosk-mode.service
[Unit]
Description=Ensure Kiosk Mode

[Service]
ExecStart=/usr/local/bin/ensure-kiosk-mode.sh

[Install]
WantedBy=multi-user.target
